# packages/systemd/clawlegion.service
[Unit]
Description=ClawLegion Multi-Agent System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/clawlegion start --with-api
WorkingDirectory=/var/lib/clawlegion
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
